# kcontrol/autostart/autostart.desktop
[Desktop Entry]
Encoding=UTF-8
Icon=launch
Type=Service
ServiceTypes=KCModule
Exec=kcmshell autostart

X-KDE-Library=autostart
X-KDE-FactoryName=autostart
X-KDE-ParentApp=kcontrol

Name=Autostart Manager

Comment=A KControl tool for managing what programs start up with KDE.

Keywords=autostart manager

Categories=Qt;KDE;X-KDE-settings-components;
